# woof-code/rootfs-skeleton/root/.config/autostart/barrier_setup.desktop
[Desktop Entry]
Type=Application
Name=Barrier
Exec=/root/my-applications/bin/barrier_setup.sh
Terminal=true
